# com.ying.myagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ying.myagent</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ying/Documents/MyAgent/.venv/bin/python</string>
        <string>/Users/ying/Documents/MyAgent/run.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/ying/Documents/MyAgent</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/ying/Documents/MyAgent/logs/myagent.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ying/Documents/MyAgent/logs/myagent.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
