# deploy/launchd/com.ying.myagent.s1check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ying.myagent.s1check</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/ying/Documents/MyAgent && bash scripts/s1_daily_check.sh >> logs/s1-daily-log.md 2>&1</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>55</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/ying/Documents/MyAgent/logs/s1check.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ying/Documents/MyAgent/logs/s1check.err.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
